# filemanager-httpd.socket
[Socket]
ListenStream=80
Accept=true
Backlog=5

[Install]
WantedBy=sockets.target
